# deploy/systemd/risk-monitor.service
[Unit]
Description=Polymarket Risk Monitor Sync
After=network.target

[Service]
Type=oneshot
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-analytics
Environment="PATH=/opt/polymarket-analytics/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/opt/polymarket-analytics/.env
ExecStart=/opt/polymarket-analytics/venv/bin/python scripts/risk_monitor.py sync
StandardOutput=append:/var/log/polymarket/risk-monitor.log
StandardError=append:/var/log/polymarket/risk-monitor.log
